# monitoring.timer.template
[Unit]
Description=Timer for test unit

[Timer]
OnUnitActiveSec=20s
OnBootSec=10s
AccuracySec=1s

[Install]
WantedBy=timers.target
