# monitoring.timer.template
[Unit]
Description=Timer for test unit

[Timer]
OnUnitActiveSec=5min
OnBootSec=10s
AccuracySec=1s

[Install]
WantedBy=timers.target
